# dev/docker/sapi/fpm/supervisord.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=1
loglevel=info

[unix_http_server]
file=/var/run/supervisord.sock
chmod=0700
username=charcoal
password=charcoal

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock
username=charcoal
password=charcoal

; Start PHP-FPM first
[program:php-fpm]
command=php-fpm -F
priority=10
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Then Nginx, using the rendered config under /home/charcoal/nginx/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /home/charcoal/nginx/nginx.conf
priority=20
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
